# stow/macos/Library/LaunchAgents/ajf.day-night-notifier.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ajf.day-night-notifier</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd $HOME ; .local/bin/macos/day-night-notifier .local/bin/common-dotfiles/day-night-update</string>
    </array>
    <key>StandardErrorPath</key>
    <string>/tmp/day-night-notifier.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/day-night-notifier.log</string>
    <key>Nice</key>
    <integer>10</integer>
  </dict>
</plist>
